# pi_zero/pizero-server.service
[Unit]
Description=Pi Zero File Server (ESP32 Node-1 compatible)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi_zero
ExecStart=/usr/bin/python3 /home/pi/pi_zero/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
